# update_edge.service
[Unit]
Description=Update Edge

[Service]
Type=simple
ExecStart=python3 /opt/eyeflow/install/update_edge.py

[Install]
WantedBy=multi-user.target
